# debian/prerm
systemctl stop amidiminder
systemctl disable amidiminder
